# files/etc/systemd/system/shop.service
[Unit]
Description=Shop service

[Service]
ExecStart=/usr/bin/shop
User=root
Group=root
Restart=always

[Install]
WantedBy=multi-user.target